# scripts/renew-cert.timer
[Unit]
Description=Renew the certificate monthly

[Timer]
OnCalendar=monthly
Unit=renew-cert.service

[Install]
WantedBy=timers.target
